# templates/vector.service.j2
[Unit]
Description=Vector
After=network-online.target
Requires=network-online.target

[Service]
User=root
Group=root
ExecStart=/usr/local/bin/vector --config /etc/vector/vector.toml
Restart=always

[Install]
WantedBy=multi-user.target

